# supervisor/laravel.ini
[supervisord]
nodaemon=true
[program:horizon]
process_name=%(program_name)s_%(process_num)02d
command=/var/www/html/artisan horizon
stopwaitsecs=3600
autostart=true
autorestart=true
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/supervisor/horizon.log
[program:octane]
process_name=%(program_name)s_%(process_num)02d
command=/var/www/html/artisan octane:start --workers=%(ENV_NUMPROCS) --task-workers=%(ENV_WORKERS) --max-requests=250
autostart=true
autorestart=true
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/supervisor/octane.log
[program:cron]
command=/usr/sbin/crond -f
autostart=true
autorestart=true
redirect_stderr=true
